# assets/org.concourse-ci.darwin-worker.plist.erb
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>AbandonProcessGroup</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>Label</key>
    <string>org.concourse-ci.darwin-worker</string>
    <key>Nice</key>
    <integer>0</integer>
    <key>ProgramArguments</key>
    <array>
    <string>/usr/local/bin/concourse</string>
    <string>worker</string>
    <string>--work-dir=<%=ENV['CONCOURSE_DIR']%>/work</string>
    <string>--name=<%=ENV['HOST']%></string>
    <string>--tsa-host=<%=ENV['TSA_HOST']%></string>
    <string>--tsa-public-key=<%=ENV['CONCOURSE_DIR']%>/config/config/tsa.pub</string>
    <string>--tsa-worker-private-key=<%=ENV['CONCOURSE_DIR']%>/config/worker.key</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
